# roles/archlinux-mediaserver/files/sc-suspend.service
[Unit]
Description=Suspend

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl suspend
