# roles/archlinux-pc/files/btrfs-snapshot.timer
[Unit]
Description=btrfs snapshot

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
