# roles/network_monitor/files/ping-healthchecks-io.timer
[Unit]
Description=Ping healthchecks.io

[Timer]
OnCalendar=*:0/4

[Install]
WantedBy=multi-user.target
