# roles/archlinux-mediaserver/templates/backup-photos.timer.j2
[Unit]
Description=Back up photos

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
