# roles/archlinux-mediaserver/templates/blocklist-update.service.j2
[Unit]
Description=Suspend

[Service]
Type=oneshot
ExecStart=/usr/bin/transmission-remote --blocklist-update
