# roles/archlinux-mediaserver/templates/blocklist-update.timer.j2
[Unit]
Description=Update BT blocklists

[Timer]
OnCalendar=*-*-* 19:00:00

[Install]
WantedBy=multi-user.target
