# roles/archlinux-pc/templates/craig-backup.timer.j2
[Unit]
Description=Back up home directory

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
